# origin-ui.service
# Copy to /etc/systemd/system/origin-ui.service
# Edit User, WorkingDirectory, and ExecStart paths to match your setup
[Unit]
Description=Origin UI
After=network.target

[Service]
Type=simple
# CHANGE THESE to match your user and install path:
User=YOURUSER
WorkingDirectory=/home/YOURUSER/origin-ui
ExecStart=/home/YOURUSER/origin-ui/venv/bin/uvicorn app:app --port 8000 --host 0.0.0.0
Restart=always
RestartSec=3
EnvironmentFile=-/home/YOURUSER/origin-ui/.env

[Install]
WantedBy=multi-user.target
